# minecraft/files/backup-server.timer
[Unit]
Description=Backup world
Requires=backup-server.service

[Timer]
OnCalendar=*-*-* 03:15:00

[Install]
WantedBy=timers.target
