# minecraft/files/save-server.timer
[Unit]
Description=Save World to Disk
Requires=save-server.service

[Timer]
OnCalendar=*-*-* *:00/30:00

[Install]
WantedBy=timers.target
